# infra/ubuntu/afrows-agent.service.sample
[Unit]
Description=Afrows Monitoring Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/afrows
EnvironmentFile=/etc/afrows/agent.env
ExecStart=/usr/bin/python3 apps/agent/run.py
Restart=on-failure
RestartSec=5
User=afrows-agent
Group=afrows-agent
StateDirectory=afrows-agent
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/lib/afrows-agent
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
