# infra/ubuntu/afrows-backend.service.sample
[Unit]
Description=Afrows Backend API
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/afrows/apps/backend
Environment=NODE_ENV=production
EnvironmentFile=/etc/afrows/afrows.env
ExecStart=/usr/bin/node dist/main.js
Restart=on-failure
RestartSec=5
User=afrows
Group=afrows
StateDirectory=afrows
LogsDirectory=afrows
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/lib/afrows /var/log/afrows
CapabilityBoundingSet=
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictSUIDSGID=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
SystemCallFilter=@system-service

[Install]
WantedBy=multi-user.target
